# scripts/monitoring/portfolio-collect-logs.cron
# /etc/cron.d/portfolio-collect-logs
# Daily Docker log collection for AgentLog analysis system.
# Runs at midnight UTC as the lukasz user.
# Output is appended to /var/log/portfolio/collect-logs.log for debugging.
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

0 0 * * * lukasz DOCKER_LOGS_DIR=/var/log/portfolio/docker-logs COMPOSE_FILE=/home/lukasz/portfolio/docker-compose.prod.yml /home/lukasz/portfolio/scripts/monitoring/collect-logs.sh >> /var/log/portfolio/collect-logs.log 2>&1
